# laravel-websockets.conf
[program:laravel-websockets]
command=php /var/www/html/Cabsisem/artisan websockets:serve --port=6001
directory=/var/www/html/Cabsisem
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/laravel-websockets.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/laravel-websockets-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stopwaitsecs=3600
killasgroup=true
stopsignal=QUIT

; Configuración adicional para producción
environment=LARAVEL_ENV="production"
startretries=3
startsecs=10

; Configuración de logs
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
stderr_capture_maxbytes=1MB
stderr_events_enabled=false